# scripts/nvr-fleet-agent.service
[Unit]
Description=NVR Fleet Agent
After=network-online.target go2rtc.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
EnvironmentFile=/etc/nvr-fleet-agent.env
ExecStart=/opt/nvr-fleet-agent/.venv/bin/python /opt/nvr-fleet-agent/agent.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
